# root/etc/services.d/crashplan/run
#!/usr/bin/with-contenv bash
. /app/install.vars

exec \
	s6-setuidgid cont 