# root/etc/services.d/crashplan/run
#!/usr/bin/with-contenv bash
exec /usr/bin/java $SRV_JAVA_OPTS -classpath $FULL_CP com.backup42.service.CPService