# root/etc/services.d/crashplan/run
#!/usr/bin/with-contenv bash
cd /app || exit 1
exec /opt/jre/bin/java $SRV_JAVA_OPTS -classpath $FULL_CP com.backup42.service.CPService